# gps/android-gpsd.service
# android-gpsd.service
#
# Systemd service to start gpsd on the Raspberry Pi and listen for GPS data
# forwarded from the Android-GPSd-Forwarder app over UDP.
#
# Android app settings:
#   Host: relectric.local
#   Port: 29998
#
# Install on the Pi:
#   sudo cp gps/android-gpsd.service /etc/systemd/system/android-gpsd.service
#   sudo systemctl disable --now gpsd.socket gpsd
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now android-gpsd.service
#
# Useful commands:
#   systemctl status android-gpsd.service
#   journalctl -u android-gpsd.service -f

[Unit]
Description=Android GPS Forwarder listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/gpsd -Nn udp://*:29998
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target